# airootfs/etc/skel/.config/autostart/appletneotheme.desktop
[Desktop Entry]
Name=AppletneoTheme
Exec=appletneotheme
Icon=utilities-terminal
Type=Application
Categories=System;Utility;